# addon/goprobe.service
[Unit]
Description=Network Traffic Monitoring
After=network-online.target syslog.target

[Service]
Type=simple
ExecStart=/usr/local/bin/goProbe -config /etc/goprobe.conf
StandardOutput=syslog
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
Alias=goprobe.service
